# deploy/supervisor/birdie-web.conf
[program:birdie-web]
command=/bin/sh -c "sleep infinity"
autostart=true
autorestart=true
stderr_logfile=/var/log/birdie-web.err.log
stdout_logfile=/var/log/birdie-web.out.log
user=root
